# infrastructure/loadtesting/terraform/shared/elasticsearch_ansible/roles/ansible/templates/ansible.timer
[Unit]
Description=Runs ansible

[Timer]
OnCalendar=*-*-* *:00:00

[Install]
WantedBy=default.target
