# tools/mac/root/Library/LaunchDaemons/co.kolide.osquery.enroll.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>co.kolide.osquery.enroll</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/osqueryd</string>
		<string>--flagfile=/etc/osquery/kolide.flags</string>
	</array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/var/log/osquery/osquery-error.log</string>
	<key>StandardOutPath</key>
    <string>/var/log/osquery/osquery-output.log</string>
</dict>
</plist>

