# schema/tables/crontab.yml
name: crontab
examples: |-
  List commands scheduled for execution as cron jobs
  
  ```
  SELECT * FROM crontab;
  ```
columns:
  - name: pid_with_namespace
    platforms:
      - windows
