# schema/tables/startup_items.yml
name: startup_items
examples: |-
  List commands executed as user/logon startup items.
  
  ```
  SELECT name, type FROM startup_items WHERE status='enabled';
  ```
