# tools/seed/fleet-mobile-policies.yml
---
# Policies for the Personal mobile devices fleet (personal laptops / BYOD)
#
# Usage:
#   fleetctl apply -f tools/seed/fleet-mobile-policies.yml --policies-fleet "Personal mobile devices"
#
# These are unique to this fleet — not duplicates of global policies.
apiVersion: v1
kind: policy
spec:
  name: Suspicious autostart (Windows)
  query: SELECT 1 WHERE NOT EXISTS (SELECT 1 FROM startup_items WHERE path = "regsvr32" AND args LIKE "%http%");
  description: "Checks for an autostart that is attempting to load a dynamic link library (DLL) from the internet."
  resolution: "Remove the suspicious startup entry."
  platform: windows
---
apiVersion: v1
kind: policy
spec:
  name: SMBv1 client driver disabled (Windows)
  query: SELECT 1 FROM windows_optional_features WHERE name = 'SMB1Protocol-Client' AND state != 1;
  description: "Checks that the SMBv1 client is disabled."
  resolution: "Contact your IT administrator to discuss disabling SMBv1 on your system."
  platform: windows
---
apiVersion: v1
kind: policy
spec:
  name: Secure keyboard entry for Terminal application enabled (macOS)
  query: SELECT 1 FROM managed_policies WHERE domain = 'com.apple.Terminal' AND name = 'SecureKeyboardEntry' AND value = 1 LIMIT 1;
  description: "Checks that a mobile device management (MDM) solution configures the Mac to enable secure keyboard entry for the Terminal application."
  resolution: "Contact your IT administrator to ensure your Mac is receiving a profile that enables secure keyboard entry for the Terminal application."
  platform: darwin
---
apiVersion: v1
kind: policy
spec:
  name: Ad tracking is limited (macOS)
  query: SELECT 1 FROM managed_policies WHERE domain='com.apple.AdLib' AND name='forceLimitAdTracking' AND value='1' LIMIT 1;
  description: "Checks that a mobile device management (MDM) solution configures the Mac to limit advertisement tracking."
  resolution: "Contact your IT administrator to ensure your Mac is receiving a profile that disables advertisement tracking."
  platform: darwin
---
apiVersion: v1
kind: policy
spec:
  name: Link-Local Multicast Name Resolution (LLMNR) disabled (Windows)
  query: SELECT 1 FROM registry WHERE path LIKE 'HKEY_LOCAL_MACHINE\SOFTWARE\Policies\Microsoft\Windows NT\DNSClient\EnableMulticast' AND CAST(data as integer) = 0;
  description: "Checks if a Group Policy configures the computer to disable LLMNR. Disabling LLMNR can prevent malicious actors from gaining access to the computer's credentials."
  resolution: "Contact your IT administrator to ensure your computer is receiving a Group Policy that disables LLMNR on your system."
  platform: windows
